# supervisor.apps.ini
[program:gunicorn]
directory = /main
command = gunicorn
autostart = true
autorestart = false
startretries = 0

[program:celery]
directory = /main
command = celery -A make_celery worker
autostart = true
autorestart = false
startretries = 0
